# snippet/crontab
10 10 * * * cd /root/vps/snippet/tg_bot; python3 telethon_send_message.py >> tg.log 2>&1

# 14 19 * * *  sh /usr/sbin/rtcwake -m mem -s 18000 >> /home/jue/SuspendResume.log 2>&1
